# daoliyu-subsonic.service
[Unit]
Description=道理鱼 Subsonic API Bridge
Documentation=https://github.com/andychao2024/daoliyu-subsonic-bridge
After=network-online.target daoliyu.music.service
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /path/to/daoliyu_subsonic_bridge.py --host 0.0.0.0 --port 4040
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=daoliyu-subsonic

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/path/to/music /path/to/daoliyu.db

[Install]
WantedBy=multi-user.target
